# debian/cron.daily
#!/bin/bash

. /etc/default/chicanery

if [ ! -z $RUN] && [ $RUN == "yes"] ; then
    echo do it
fi

